# minecraft-java.socket
[Socket]
ListenFIFO=%t/minecraft-java.stdin
Service=minecraft-java.service

[Install]
WantedBy=multi-user.target
